# userspace/files/adbd.service
[Unit]
Description=adbd Service

[Service]
User=root
ExecStart=/usr/kommu/adbd
