# userspace/files/gpio.service
[Unit]
Description=GPIO export

[Service]
Type=oneshot
ExecStart=/usr/kommu/gpio.sh

[Install]
WantedBy=multi-user.target
