# userspace/files/wifi.service
[Unit]
Description=Wifi

[Service]
Type=oneshot
ExecStart=/sbin/insmod /usr/kommu/wlan.ko

[Install]
WantedBy=multi-user.target
